# package/base-files/files/etc/uci-defaults/99-auto-extroot
#!/bin/sh
# Auto-extroot script: try to find a suitable partition to use as /overlay
# This script will run once on first boot (uci-defaults)
LOG=/tmp/extroot.log
echo "[extroot] start" > ${LOG}

/sbin/block info > /tmp/blockinfo.txt 2>&1 || true

# Find a block device with available space and not mounted
for dev in $(ls /dev | egrep 'sd[b-z][1-9]?|mmcblk[0-9]p[1-9]' | sort); do
  devpath="/dev/${dev}"
  # skip if partition is mounted
  if mount | grep -q "${devpath}"; then
    continue
  fi
  # check size (use lsblk to get size in bytes)
  size=$(lsblk -b -dn -o SIZE ${devpath} 2>/dev/null || echo 0)
  if [ -z "${size}" ]; then size=0; fi
  # require at least 1GB free
  if [ ${size} -ge 1073741824 ]; then
    echo "[extroot] candidate ${devpath} size=${size}" >> ${LOG}
    # format as ext4 and mount to /overlay
    mkfs.ext4 -F ${devpath} >> ${LOG} 2>&1 || { echo "[extroot] format failed" >> ${LOG}; continue; }
    mkdir -p /mnt/newoverlay
    mount ${devpath} /mnt/newoverlay || { echo "[extroot] mount failed" >> ${LOG}; continue; }
    # copy current overlay content
    tar -C /overlay -cpf - . 2>/dev/null | tar -C /mnt/newoverlay -xpf - 2>/dev/null || true
    umount /mnt/newoverlay || true
    # write fstab to mount as overlay on boot
    echo "${devpath} /overlay ext4 defaults 0 0" >> /etc/fstab
    echo "[extroot] success: ${devpath}" >> ${LOG}
    break
  fi
done

echo "[extroot] finished" >> ${LOG}
exit 0
